# crontab
@weekly /home/pi/OnionPi/get_ad_list.sh

